# funretro.service
[Unit]
Description=FunRetro

[Service]
User=funretro
WorkingDirectory=/opt/funretro
ExecStart=/usr/bin/gulp

[Install]
WantedBy=multi-user.target
