# funretro.service
[Unit]
Description=FunRetro

[Service]
WorkingDirectory=/vagrant
ExecStart=/usr/bin/gulp

[Install]
WantedBy=multi-user.target
